# tests/vm_templates/linux/cdrom_autoexec.service
[Unit]
Description=Mount and execute CDROMs

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/cdrom_autoexec.sh