# tests/vm_templates/linux/cdrom.rules
SUBSYSTEM=="block", ENV{ID_CDROM}=="?*", ACTION=="change", RUN+="/bin/systemctl --no-block start cdrom_autoexec.service"
