# tests/vm_templates/linux/cdrom.rules
SUBSYSTEM=="block", ENV{ID_CDROM}=="?*", ACTION=="change", RUN+="/usr/local/sbin/cdrom_autoexec.sh"
