# tests/vm_templates/linux/rc.local
#!/bin/sh -e

dmesg -n 8
echo 12000 > /sys/block/sr0/device/timeout

exit 0
